# update_server/deploy/amap-companion-sync.timer
[Unit]
Description=Periodically sync and build AMap Companion APK

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=amap-companion-sync.service
Persistent=true

[Install]
WantedBy=timers.target
